# hwctld.service
[Unit]
Description=Daemon for controlling various system devices
[Service]
Type=exec
ExecStart=/usr/bin/hwctld
[Install]
WantedBy=multi-user.target